# systemd/system/getty@.service.d/finnix.conf
[Service]
ExecStart=
#ExecStart=-/sbin/agetty --noclear --autologin root --noissue --skip-login %I $TERM
ExecStart=-/sbin/shell-login %I
TTYVTDisallocate=no
